# systemd/flashpasted.service
[Unit]
Description=flashpaste daemon (clipboard owner + paste dispatcher)
Documentation=https://github.com/NagyVikt/flashpaste
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/flashpasted
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
